# packages/worker/transcription-service-worker.service
[Unit]
Description=Worker app for the guardian's transcription service
Documentation=https://github.com/guardian/transcription-service
After=network.target

[Service]
Type=simple
User=ubuntu
ExecStart=/usr/bin/node /opt/transcription-service/index.js
Restart=on-failure
Environment="LD_LIBRARY_PATH=/usr/local/cuda-12.4/lib:/usr/local/cuda-12.4/lib64"

[Install]
WantedBy=multi-user.target
